# production/config_v2.0-ms/octal_app.ini
[program:octal]
command = /srv/run_octal_server.sh    ;command to start app
user = octal
stdout_logfile = /srv/log/octal.log
redirect_stderr = true
autostart = true
autorestart = true
