# production/config_v1.1-pilot/octal_content.ini
[program:octal_content]
command = /srv/run_content_server.sh
user = octal
stdout_logfile = /srv/log/content_server.log
redirect_stderr = true
autostart = true
autorestart = true
