# production/config_v1.1-pilot/octal_app.ini
[program:octal_app]
command = /srv/run_app_server.sh    ;command to start app
user = octal
stdout_logfile = /srv/log/app_server.log
redirect_stderr = true
autostart = true
autorestart = true
